# oldBird/rmi/rmi-provider/src/main/resources/rmi.conf
[program:rmi]
directory = /data/app/;
command = java WarehouseServer;
autostart = true ;
autorestart = false ;
user = root ;